# service.ini
[program:service]
command=<App start command>
directory=<app root home folder>
autorestart=true
startretries=300
priority=150
stdout_logfile=/dev/stdout
stdout_logfile_maxbytes=0
stderr_logfile=/dev/stderr
stderr_logfile_maxbytes=0


[program:rsyslogd]
command=/bin/ash -c "until [ -f /etc/rsyslog.d/00-elasticsearch.conf ]; do echo 'not found'; sleep 5; done; echo 'found'; kill -9 $(cat /var/run/rsyslogd.pid); rsyslogd"
startsecs=0
autorestart=false
startretries=300
priority=999
stdout_logfile=/dev/stdout
stdout_logfile_maxbytes=0
stderr_logfile=/dev/stderr
stderr_logfile_maxbytes=0

